# MicroserviceTest/supervisord.conf
[supervisord]
nodaemon=true
 
[program:filebeat]
command=filebeat -e -c=/etc/filebeat/filebeat.yml

